# novnc/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -shared -clip 566x406+45+45
autorestart=true
